# bin/crossword.service
[Unit]
Description="crossword server"
Wants=network-online.target
After=network-online.target

[Service]
Type=Simple
ExecStart=/apps/crossword-factory/env/bin/gunicorn -c /apps/crossword-factory/bin/gunicorn.conf.py -b unix:/apps/crossword-factory/var/django.sock
Restart=always
Environment="PATH=/apps/crossword-factory/env/bin/"
WorkingDirectory=/apps/crossword-factory/
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutSec=6
PrivateTmp=true

[Install]
WantedBy=default.target
